# etc/services/fpkeeper/run
#!/bin/sh
exec 2>&1
echo >&2 "Starting fpkeeper"
export PREDATOR_HOME="/opt/predator"
exec "$PREDATOR_HOME"/installed/storekeeper/apps/fpkeeper/fpkeeper.py
