# etc/services/agentfront/run
#!/bin/sh
exec 2>&1
echo >&2 "Starting agentfront"
export PREDATOR_HOME="/opt/predator"
exec "$PREDATOR_HOME"/installed/agent/apps/agentfront/agent.py
