# etc/services/sanitizer/run
#!/bin/sh
exec 2>&1
echo >&2 "Starting sanitizer"
export PREDATOR_HOME="/opt/predator"
exec "$PREDATOR_HOME"/installed/preprocessor/apps/sanitizer/sanitizer.py
